# supervisord.conf
[supervisord]
nodaemon=true

[program:twemproxy]
command=/twemproxy -c /etc/conf/twemproxy.yml
numprocs=1
autostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true



